# modelix-bot.service
[Unit]
Description=Modelix Telegram Notification Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/modelix-bot
# TELEGRAM_BOT_TOKEN, TELEGRAM_CHANNEL_ID, TELEGRAM_PROXY_URL, опционально пути и т.д.
EnvironmentFile=-/var/www/modelix-bot/.env
ExecStart=/var/www/modelix-bot/venv/bin/python /var/www/modelix-bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
